# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# =============================================================================
# Virtual Display Service
# =============================================================================

[program:xorg]
command=/usr/lib/xorg/Xorg :99 -config /etc/xorg.conf -noreset -nolisten tcp
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xorg.log
stderr_logfile=/var/log/supervisor/xorg-error.log
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

# =============================================================================
# Audio Service
# =============================================================================

[program:pulseaudio]
command=/usr/bin/pulseaudio --system --exit-idle-time=-1 -nF --log-target=stderr --disable-shm=yes --disallow-exit
autostart=true
autorestart=true
priority=105
stdout_logfile=/var/log/supervisor/pulseaudio.log
stderr_logfile=/var/log/supervisor/pulseaudio-error.log
environment=XDG_RUNTIME_DIR=/run/user/0

# =============================================================================
# Window Manager on :99 (required for Chrome/GoExport)
# =============================================================================

[program:openbox99]
command=/usr/bin/openbox
environment=DISPLAY=:99
autostart=true
autorestart=true
priority=110
stdout_logfile=/var/log/supervisor/openbox99.log
stderr_logfile=/var/log/supervisor/openbox99-error.log

# =============================================================================
# VNC Desktop Services (optional - for remote desktop access)
# =============================================================================

[program:x11vnc99]
command=/usr/bin/x11vnc -display :99 -rfbport 5999 -nopw -nocursorshape -forever -noxfixes -noxdamage -bg
environment=DISPLAY=:99,XAUTHORITY=/root/.Xauthority
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/x11vnc99.log
stderr_logfile=/var/log/supervisor/x11vnc99_error.log

[program:tigervnc]
command=/usr/bin/vncserver :1 -geometry 1920x1080 -depth 24 -SecurityTypes None -fg
autostart=true
autorestart=true
priority=101
stdout_logfile=/var/log/supervisor/tigervnc.log
stderr_logfile=/var/log/supervisor/tigervnc_error.log

[program:novnc99]
command=/usr/bin/websockify --web /usr/share/novnc 6099 localhost:5999
autostart=true
autorestart=true
priority=102
stdout_logfile=/var/log/supervisor/novnc99.log
stderr_logfile=/var/log/supervisor/novnc99_error.log

[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5901
autostart=true
autorestart=true
priority=102
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_error.log

# =============================================================================
# Web Services
# =============================================================================

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autostart=true
autorestart=true
priority=210
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm-error.log

# =============================================================================
# Laravel Queue Workers
# =============================================================================

[program:laravel-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
priority=300
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-default.log
stopwaitsecs=3600

[program:laravel-worker-exports]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=exports --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
priority=310
environment=DISPLAY=:99,XAUTHORITY=/root/.Xauthority,DBUS_SYSTEM_BUS_ADDRESS="unix:path=/var/run/dbus/system_bus_socket"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-worker-exports.log
stopwaitsecs=3600
startsecs=15

# =============================================================================
# Laravel Scheduler (Cron)
# =============================================================================

[program:laravel-scheduler]
command=/bin/bash -c "while [ true ]; do php /var/www/html/artisan schedule:run --verbose --no-interaction; sleep 60; done"
autostart=true
autorestart=true
user=www-data
priority=320
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log

# =============================================================================
# Additional Services (Extensible Section)
# =============================================================================
# Add your custom services below this line
# Example:
#
# [program:custom-service]
# command=/path/to/your/command
# autostart=true
# autorestart=true
# user=www-data
# priority=400
# stdout_logfile=/var/log/supervisor/custom-service.log
# stderr_logfile=/var/log/supervisor/custom-service-error.log
